# pilab/deploy/pilab-api.service
[Unit]
Description=PiLab Central Data API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/pilab
Environment=PILAB_DB_PATH=/opt/pilab/data/pilab.db
ExecStart=/opt/pilab/.venv/bin/uvicorn pilab.api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
