# pilab/deploy/pilab-scout.service
[Unit]
Description=PiLab Scout Agent
After=network.target

[Service]
Type=simple
User=mypi
WorkingDirectory=/home/mypi/Projects/multiagent_researchlab
Environment=PILAB_API_BASE_URL=http://mypi:8000
Environment=PILAB_LLM_SCOUT_URL=http://localhost:8080
ExecStart=/home/mypi/Projects/multiagent_researchlab/.venv/bin/python -m pilab.scout.agent
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
