# pilab/deploy/pilab-api.service
[Unit]
Description=PiLab Central Data API
After=network.target

[Service]
Type=simple
User=mypi
WorkingDirectory=/home/mypi/Projects/multiagent_researchlab
Environment=PILAB_DB_PATH=/mnt/pilab/pilab.db
ExecStart=/home/mypi/Projects/multiagent_researchlab/.venv/bin/uvicorn pilab.api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
