# pilab/deploy/pilab-evaluator.service
[Unit]
Description=PiLab Evaluator Agent
After=network.target pilab-api.service

[Service]
Type=simple
WorkingDirectory=/opt/pilab
Environment=PILAB_LLM_EVALUATOR_MODEL=qwen2.5:14b
ExecStart=/opt/pilab/.venv/bin/python -m pilab.evaluator.agent
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
