# pilab/deploy/pilab-scout.service
[Unit]
Description=PiLab Scout Agent
After=network.target pilab-api.service

[Service]
Type=simple
WorkingDirectory=/opt/pilab
Environment=PILAB_LLM_SCOUT_MODEL=phi3:mini
ExecStart=/opt/pilab/.venv/bin/python -m pilab.scout.agent
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
