# pilab/deploy/pilab-planner.service
[Unit]
Description=PiLab Planner Agent
After=network.target pilab-api.service

[Service]
Type=simple
WorkingDirectory=/opt/pilab
Environment=PILAB_LLM_PLANNER_MODEL=mistral:7b
ExecStart=/opt/pilab/.venv/bin/python -m pilab.planner.agent
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
