# pilab/deploy/pilab-planner.service
[Unit]
Description=PiLab Planner Agent
After=network.target

[Service]
Type=simple
User=mypi3
WorkingDirectory=/home/mypi3/Projects/multiagent_researchlab
Environment=PILAB_API_BASE_URL=http://192.168.1.35:8000
Environment=PILAB_LLM_PLANNER_URL=http://localhost:8082
ExecStart=/home/mypi3/Projects/multiagent_researchlab/.venv/bin/python -m pilab.planner.agent
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
